# kafka2influxdb.upstart
description "kafka2influxdb"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

console log
setuid _kafka2influxdb
setgid _kafka2influxdb
kill signal SIGTERM
kill timeout 10

script
	. /etc/default/kafka2influxdb
	exec /usr/local/bin/kafka2influxdb start --logfile=$LOGFILE --loglevel=$LOGLEVEL
end script

